# digest/crontab
0 12 * * * /usr/local/bin/python /app/main.py
